# docs/discord-bot.service.example
[Unit]
Description=PanzerVault Bot
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/panzervault-bot
EnvironmentFile=/opt/panzervault-bot/.env
ExecStart=/usr/bin/node /opt/panzervault-bot/src/index.js
Restart=always
RestartSec=10
User=discordbot
Group=discordbot

[Install]
WantedBy=multi-user.target
